# rootfs/etc/systemd/system/expand-rootfs.service
[Unit]
Description=Expand root filesystem to fill disk
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/expand-rootfs.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
